# system-setup/lenovo-IdeaPad-5-15ARE05-setup/touchpadfix.service
[Unit]
Description=Fix touchpad issue by binding correct driver

[Service]
ExecStart=/usr/local/bin/touchpadfix
Type=oneshot
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
